# frame/systemd/birdframe.service
[Unit]
Description=AvianVisitors e-ink frame update
Documentation=https://github.com/Twarner491/AvianVisitors
Wants=network-online.target
After=network-online.target

[Service]
# install.sh rewrites User= and the paths to match your install.
Type=oneshot
User=monalisa
WorkingDirectory=/home/monalisa/AvianVisitors/frame
ExecStart=/home/monalisa/AvianVisitors/frame/.venv/bin/python /home/monalisa/AvianVisitors/frame/display.py --config /home/monalisa/.birdframe/config.toml
Environment=PYTHONUNBUFFERED=1
Nice=10
TimeoutStartSec=180
