# avian/forwarding/avian-mqtt.service
[Unit]
Description=AvianVisitors MQTT bridge
After=network-online.target

[Service]
Type=simple
# Edit User= and the ExecStart path to match your Pi's username if
# you didn't install BirdNET-Pi as the `birdnet` user.
User=birdnet
ExecStart=/usr/bin/python3 %h/avian-mqtt.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
